# etc/udev/rules.d/99-monitor_hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/systemctl start monitor_hotplug.service"
